# code/services/webcast-monitor.timer
[Unit]
Description=Check for changed webcast

[Timer]
OnStartupSec=1s
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
