# code/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=443
